# grub-monitor.path
[Unit]
Description=Watch grub.cfg for changes

[Path]
PathChanged=/boot/grub/grub.cfg

[Install]
WantedBy=multi-user.target